# target/linux/ar71xx/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

if ! [ -e /etc/config/network ]
then
	echo "Network config not found"
	exit 1
fi

changed="0"

#Get saved  3G and LAN settings
list=`/sbin/keep_settings.sh get sim1 sim2 lan ppp`
if [ -n "$list" ]
then
	echo "Restoring saved 3G and LAN settings"
	for item in $list; do
		uci set $item
	done

	ifname=`uci -q get network.ppp.ifname`

	if [ $ifname ]; then
		uci -q set network.wan.ifname=$ifname
		if [ "$ifname" == "eth2" ]; then
			uci set network.wan.proto="dhcp"
		elif [ "$ifname" == "3g-ppp" ] || [ "$ifname" == "wwan0" ]; then
			uci set network.wan.proto="none"
		fi
	fi

	/sbin/keep_settings.sh delete
	echo "$list" > /tmp/temp_config
	changed="1"
	#Don't do firstboot stuff
	uci -q delete teltonika.sys.first_login=0
	uci commit teltonika
fi
#Get saved SIM PIN
pin1=$(/sbin/mnf_info simpin)
pin2=$(/sbin/mnf_info sim2pin)
[ -n "$pin1" ] && [ "$pin1" != "N/A" ] && {
	echo "Restoring saved SIM1 PIN"
	uci set network.ppp=interface
	uci set network.ppp.pincode="$pin1"
	uci set simcard.sim1.pincode="$pin1"
	changed="1"
}

[ -n "$pin2" ] && [ "$pin2" != "N/A" ] && {
	echo "Restoring saved SIM2 PIN"
	uci set simcard.sim2.pincode="$pin2"
	changed="1"
}


[ "$changed" = "1" ] && {
	uci commit simcard
	uci commit network
}

exit 0

