# custom_feeds/teltonika_luci/i18n/turkish/root/etc/uci-defaults/luci-i18n-turkish
#!/bin/sh
uci batch <<-EOF
	set luci.languages.tr=Turkish
        commit luci
EOF

