# package/teltonika/mobutils/bin/files/sms_utils.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 Teltonika

START=97
#STOP=19
#count=0

process_start()
{
	/usr/sbin/sms_utils startup &
}

boot()
{
	process_start
}

start()
{
	process_start
}
