# custom_feeds/teltonika_luci/applications/luci-cfgwzd/root/etc/uci-defaults/teltonika
config system 'sys'
	option first_login '1'
	option module 'option'
	option shw3g '0'