# package/base-files/files/etc/uci-defaults/02_uhttpd_config
#!/bin/sh
redirect_https=`uci get -q uhttpd.main.redirect_https`


if [ -n "$redirect_https" ]; then
	uci delete uhttpd.main.redirect_https
	uci commit uhttpd
fi
exit 0
