# package/base-files/files/etc/uci-defaults/04_remove-vrrp-ping-config
#!/bin/sh

uci delete vrrpd.ping

if [[ "$(uci get -q vrrpd.vid1.enabled)" != 1 ]]; then
        uci delete vrrpd.vid1
fi
