# package/base-files/files/etc/uci-defaults/12_cli_port
#!/bin/sh
cli_port=`uci get -q cli.status.port`
old_port=`echo $cli_port | grep "-"`

if [ -z "$old_port" ]; then
	uci set cli.status.port="4200-4220"
	uci commit cli
fi

exit 0
