# package/base-files/files/etc/uci-defaults/40_luci_reload_checkall
#!/bin/sh
# Copyright (C) 2015 Teltonika

luci-reload --check_all >/dev/null 2>/dev/null
