# package/base-files/files/etc/uci-defaults/83_etherwake_migrate
#!/bin/sh

. /lib/functions.sh

CONF="etherwake"

[ -f "/etc/config/$CONF" ] || return 0
uci_set "$CONF" setup interface br-lan
uci_commit "$CONF"