# custom_feeds/teltonika_luci/themes/teltonikaExp_safemode/root/etc/uci-defaults/luci-theme-teltonikaExp
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Teltonika_Experimental=/luci-static/teltonikaExp
	commit luci
EOF
exit 0
