# target/linux/ath79/teltonika/rut2_migration/uci-defaults/004-rut2-udprelay-migration.sh
#!/bin/sh

. /lib/functions.sh

[ -f "/etc/config/teltonika" ] || return 0

uci_rename udprelay br_lan lan
uci_commit udprelay

exit 0
