# target/linux/ath79/teltonika/rut2_migration/uci-defaults/007-rut2-etherwake-migration.sh
#!/bin/sh

. /lib/config/uci.sh

[ -f "/etc/config/teltonika" ] || return 0

uci_remove etherwake setup interface
uci_commit etherwake
