# target/linux/ath79/teltonika/rut2_migration/uci-defaults/117-rut2-remove-config-files.sh
#!/bin/sh

[ -f "/etc/config/teltonika" ] || return 0

rm /etc/config/auto_update /etc/config/racoon \
    /etc/config/sim_idle_protection /etc/config/sim_switch \
    /etc/config/gps