# custom_feeds/teltonika_luci/themes/default/root/etc/uci-defaults/luci-theme-default
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Default=/luci-static/default
	commit luci
EOF
exit 0
