# target/linux/ath79/teltonika/rut9_migration/uci-defaults/018-rut9-rpcd-migration.sh
#!/bin/sh

[ -f "/etc/config/teltonika" ] || return 0

cp /rom/etc/config/rpcd /etc/config/rpcd


