# target/linux/ath79/teltonika/rut9_migration/uci-defaults/037-rut9-smpp-migration.sh
#!/bin/sh

[ -e "/etc/config/teltonika" ] || return

mv /etc/config/smpp_config /etc/config/smpp
