# target/linux/ath79/teltonika/rut9_migration/uci-defaults/103-rut9-mdcollect-migration.sh
#!/bin/sh

. /lib/functions.sh

config_load mdcollectd
config_get ignore config ignore
[ -n "$ignore" ] || {
	uci_add_list mdcollectd config ignore "lo"
	uci_add_list mdcollectd config ignore "wwan0"
}

uci_commit mdcollectd
