# target/linux/ath79/generic/base-files/etc/hotplug.d/firmware/11-ath10k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath10k/cal-pci-0000:00:00.0.bin")
	case $board in
	teltonika,tcr1xx |\
	teltonika,otd1xx)
		caldata_extract "art" 0x5000 0x844
		;;
	esac
	;;
*)
	exit 1
	;;
esac
