# package/base-files/files/etc/hotplug.d/iface/02-multiwan
#!/bin/sh
[ "$ACTION" = "ifup" -o "$ACTION" = "ifupdate" ] || exit 0
killall -USR1 mwan
